# supervisord.conf
[program:redis]
directory=%(here)s
command=%(here)s/src/redis-server redis.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
autorestart=true